# rootfs/etc/systemd/system/cockpit.socket.d/listen.conf
[Socket]
ListenStream=
ListenStream=127.0.0.1:9090
ListenStream=172.18.0.1:9090
FreeBind=yes
