# rootfs/etc/systemd/system/ppp0-watchdog.service
[Unit]
Description=Ensure ppp0 interface exists; bring it up if missing
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/ppp0-watchdog.sh
